# services/ai-chatbot.service
[Unit]
Description=AI Chatbot Orchestration Service
After=network.target ollama.service
Requires=ollama.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ai-chatbot.py
Restart=on-failure
RestartSec=5s
User=root
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
